# ssh-agent/ssh-agent.service
[Unit]
Description=SSH key agent

[Service]
Type=simple
ExecStart=/usr/bin/ssh-agent -D -t 43200

[Install]
WantedBy=default.target
