# cron/crontab
59 13 * * * /home/jumski/.dotfiles/bin/backup-plasma-configs
