# todoist-helper/crontab
# Update todoist urgent tasks cache every minute
* * * * * /home/jumski/.dotfiles/todoist-helper/bin/update-urgent-cache.sh >/dev/null 2>&1